# bootc/evilcentos9/config/evilshare/usr/lib/systemd/system/smb.service.d/require_external_hdd.conf
[Unit]
RequiresMountsFor=/var/storage/external_hdd/proxmox/backups